# provision/server/.apps/desktop-control/gdm-control.service.template
[Unit]
Description=GDM Control HTTP Service
After=network.target

[Service]
Type=simple
ExecStart=__INSTALL_DIR__/gdm-control.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target

